# pro/launchd/com.example.self-heal.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!--
  self-heal nightly runner (macOS launchd). Runs the synthetic user + fix-agent every night.

  Install:
    1. Replace __SELF_HEAL_DIR__ with the absolute path to pro/self-heal.
    2. cp this file to ~/Library/LaunchAgents/com.example.self-heal.plist
    3. launchctl load ~/Library/LaunchAgents/com.example.self-heal.plist

  Note: launchd jobs often can't read external/network volumes due to macOS privacy (TCC).
  If pro/self-heal lives on such a volume, stage a copy under ~/Library and point here at that.

  Linux: use cron instead — `20 4 * * *  /path/to/pro/self-heal/nightly.sh`
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.self-heal</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__SELF_HEAL_DIR__/nightly.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>4</integer>
    <key>Minute</key><integer>20</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>__SELF_HEAL_DIR__/nightly.log</string>
  <key>StandardErrorPath</key>
  <string>__SELF_HEAL_DIR__/nightly.err.log</string>
</dict>
</plist>
